# migration/launchagents/com.haidilao.claude.log-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.haidilao.claude.log-monitor — Claude Code cron job
  Generated by: pnpm run setup
  Prompt: prompts/log-monitor.md
  Interval: 1800s (30m)
  Model: sonnet
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.haidilao.claude.log-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/hongming-claw/.local/bin/claude --print --dangerously-skip-permissions --model sonnet --add-dir /Users/hongming-claw/haidilao-business-intelligent --add-dir /Users/hongming-claw/haidilao-automation-monorepo -p "$(cat '/Users/hongming-claw/haidilao-business-intelligent/prompts/log-monitor.md')"</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hongming-claw/haidilao-automation-monorepo</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/hongming-claw/haidilao-business-intelligent/data/cron-logs/log-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hongming-claw/haidilao-business-intelligent/data/cron-logs/log-monitor.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/hongming-claw</string>
        <key>PATH</key>
        <string>/Users/hongming-claw/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>TimeOut</key>
    <integer>120</integer>
</dict>
</plist>